# uwsgi.service
[Service]
ExecStart=/home/leeps/venv/bin/uwsgi --ini /home/leeps/leeps_website/leeps_website_uwsgi.ini
KillSignal=SIGQUIT
